# deploy/lichen-faucet.service
[Unit]
Description=Lichen Faucet Service
After=lichen-validator-testnet.service
Wants=lichen-validator-testnet.service

[Service]
Type=simple
User=lichen
Group=lichen
WorkingDirectory=/var/lib/lichen

ExecStart=/usr/local/bin/lichen-faucet
Restart=always
RestartSec=5

Environment=PORT=9100
Environment=RPC_URL=http://127.0.0.1:8899
Environment=NETWORK=testnet
Environment=MAX_PER_REQUEST=10
Environment=DAILY_LIMIT_PER_IP=50
Environment=COOLDOWN_SECONDS=60
Environment=AIRDROPS_FILE=/var/lib/lichen/airdrops.json
Environment=FAUCET_KEYPAIR=/var/lib/lichen/faucet-keypair-testnet.json
Environment=RUST_LOG=info
Environment=TRUSTED_PROXY=127.0.0.1,::1

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/lichen

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lichen-faucet

[Install]
WantedBy=multi-user.target
