# deploy/lichen-custody.service
[Unit]
Description=Lichen Custody Bridge
After=network-online.target lichen-validator-testnet.service
Wants=network-online.target
Requires=lichen-validator-testnet.service

[Service]
Type=simple
User=lichen
Group=lichen
ExecStart=/usr/local/bin/lichen-custody
EnvironmentFile=/etc/lichen/custody-env
WorkingDirectory=/var/lib/lichen
Restart=always
RestartSec=5
LimitNOFILE=65535

# Security Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/lichen/custody-db
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lichen-custody

[Install]
WantedBy=multi-user.target
