# deploy/moltchain-faucet.service
[Unit]
Description=MoltChain Faucet Service
After=moltchain-validator-testnet.service
Wants=moltchain-validator-testnet.service

[Service]
Type=simple
User=moltchain
Group=moltchain
WorkingDirectory=/var/lib/moltchain

ExecStart=/usr/local/bin/moltchain-faucet
Restart=always
RestartSec=5

Environment=PORT=9100
Environment=RPC_URL=http://127.0.0.1:8899
Environment=NETWORK=testnet
Environment=MAX_PER_REQUEST=10
Environment=DAILY_LIMIT_PER_IP=50
Environment=COOLDOWN_SECONDS=60
Environment=AIRDROPS_FILE=/var/lib/moltchain/airdrops.json
Environment=FAUCET_KEYPAIR=/var/lib/moltchain/faucet-keypair-testnet.json
Environment=RUST_LOG=info
Environment=TRUSTED_PROXY=127.0.0.1,::1

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/moltchain

StandardOutput=journal
StandardError=journal
SyslogIdentifier=moltchain-faucet

[Install]
WantedBy=multi-user.target
