# deploy/moltchain-validator.service
[Unit]
Description=MoltChain Validator Node
Documentation=https://developers.moltchain.network/validator.html
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=moltchain
Group=moltchain

# Network is set via /etc/moltchain/env (MOLTCHAIN_NETWORK=testnet|mainnet)
# Port assignments: testnet → RPC=8899 WS=8900 P2P=7001
#                   mainnet → RPC=9899 WS=9900 P2P=8001
ExecStart=/usr/local/bin/moltchain-validator \
    --network ${MOLTCHAIN_NETWORK} \
    --listen-addr 0.0.0.0 \
    --rpc-port ${MOLTCHAIN_RPC_PORT} \
    --ws-port ${MOLTCHAIN_WS_PORT} \
    --p2p-port ${MOLTCHAIN_P2P_PORT} \
    --db-path /var/lib/moltchain/state-${MOLTCHAIN_NETWORK} \
    ${MOLTCHAIN_EXTRA_ARGS}

Restart=always
RestartSec=10
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/moltchain /var/log/moltchain
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Environment — MOLTCHAIN_HOME pins the P2P identity directory so cert
# files are always found at the same path regardless of HOME, db-path, or
# binary location changes.  This prevents accidental identity regeneration
# (and TOFU violations) during updates or deployments.
# HOME is also overridden for ZK keys and validator keypairs.
Environment=RUST_LOG=info
Environment=HOME=/var/lib/moltchain
Environment=MOLTCHAIN_HOME=/var/lib/moltchain
Environment=MOLTCHAIN_CONTRACTS_DIR=/var/lib/moltchain/contracts
EnvironmentFile=/etc/moltchain/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moltchain-validator

[Install]
WantedBy=multi-user.target
